# systemd/twitch-clip.service
[Unit]
Description=Twitch Clips Local Player Backend
Documentation=https://github.com/BigChiefRick/twitch-clips-local-player
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/twitch-clips-player
ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=20
KillMode=process

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=VIDEO_URL_PATH=/clips/videos

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html/twitch-clips-player/downloaded_clips

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twitch-clips

[Install]
WantedBy=multi-user.target